# tools/whatsapp-incoming-bridge/launchd/com.firstchord.whatsapp-incoming-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "https://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.firstchord.whatsapp-incoming-bridge</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>source "$HOME/.zshrc" &gt;/dev/null 2&gt;&amp;1 || true; cd "/Users/finnlemarinel/Desktop/FirstChord/music-school-dashboard/tools/whatsapp-incoming-bridge"; exec npm start</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <!-- Crash-loop guard: wait at least 30s between relaunches so a fast-exiting
       process (e.g. the watchdog firing) can't spin. -->
  <key>ThrottleInterval</key>
  <integer>30</integer>

  <key>EnvironmentVariables</key>
  <dict>
    <key>SYNC_GROUPS_ON_START</key>
    <string>true</string>
  </dict>

  <key>WorkingDirectory</key>
  <string>/Users/finnlemarinel/Desktop/FirstChord/music-school-dashboard/tools/whatsapp-incoming-bridge</string>

  <key>StandardOutPath</key>
  <string>/Users/finnlemarinel/Desktop/FirstChord/music-school-dashboard/tools/whatsapp-incoming-bridge/logs/launchagent.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/finnlemarinel/Desktop/FirstChord/music-school-dashboard/tools/whatsapp-incoming-bridge/logs/launchagent.err.log</string>
</dict>
</plist>
